# config/runit/proxmox-backup-proxy/run
#!/bin/sh
exec chpst -u root proxmox-backup-proxy 2>&1
