# config/runit/proxmox-backup-proxy/run
#!/bin/sh
exec chpst -u root /usr/bin/proxmox-backup-proxy 2>&1
