# config/runit/proxmox-backup-api/run
#!/bin/sh
exec chpst -u backup /usr/bin/proxmox-backup-api 2>&1
