# config/runit/proxmox-backup-api/run
#!/bin/sh
exec chpst -u backup proxmox-backup-api-daemon 2>&1
